# systemd/glowf1sh-cloud-client.service
[Unit]
Description=Glowf1sh Cloud Remote Control Client
Documentation=https://github.com/glowf1sh/glowboxremote
After=network-online.target glowf1sh-api-server.service
Wants=network-online.target
Requires=glowf1sh-api-server.service

[Service]
Type=simple
EnvironmentFile=/opt/glowf1sh-remote/config/cloud_api.env
ExecStart=/usr/bin/python3 /opt/glowf1sh-remote/belabox-api/cloud_client.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=glowf1sh-cloud
Restart=always
RestartSec=10

# Security
User=root
ProtectSystem=strict
ReadWritePaths=/opt/glowf1sh-remote /var/log /usr/bin /etc/systemd/system
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
