# systemd/glowf1sh-api-server.service
[Unit]
Description=Glowf1sh BelaBox REST API Server
Documentation=https://github.com/glowf1sh/glowboxremote
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/glowf1sh-remote/belabox-api/api_server.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=glowf1sh-api
Restart=always
RestartSec=10

# Security
User=root
ProtectSystem=strict
ReadWritePaths=/opt/glowf1sh-remote /opt/belaUI /tmp
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
